# .config/systemd/user/alttab.service
[Unit]
Description=alttab

[Service]
Type=simple
ExecStart=/usr/bin/alttab
Restart=always