# .config/systemd/user/cbatticon.service
[Unit]
Description=cbatticon

[Service]
Type=simple
ExecStart=/usr/bin/cbatticon
Restart=always
