# .config/systemd/user/clipster.service
[Unit]
Description=clipster daemon

[Service]
Type=simple
ExecStart=/home/maciej/projects/clipster/clipster -d
Restart=always
