# .config/systemd/user/udiskvm.service
[Unit]
Description=udisksvm

[Service]
Type=simple
ExecStart=/usr/bin/udisksvm
Restart=always
