# .config/systemd/user/volwheel.service
[Unit]
Description=Volwheel

[Service]
Type=simple
ExecStart=/usr/bin/volwheel
Restart=always